# roles/deploy-postgresql/templates/service/db_backup_full.timer.j2
[Unit]
Description=Run Full DB dump service everyday at 21:30

[Timer]
Unit=db_backup_full_{{ hospital }}.service
OnCalendar=Sun *-*-* 21:30:00

[Install]
WantedBy=timers.target